# raspberry-backend/systemd/masjid-api.service
[Unit]
Description=Masjid Display Raspberry Pi API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/opt/masjid-display/raspberry-backend
Environment=MASJID_API_KEY=
Environment=MASJID_HMAC_SECRET=
Environment=MASJID_STATE_PATH=/home/admin/masjid_state.json
Environment=MASJID_UPLOAD_DIR=/home/admin/masjid_assets/uploads
ExecStart=python3 -m uvicorn masjid_server:app --host 0.0.0.0 --port 8787
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
